# kubernetes/files/systemd/system/kube-scheduler.service.j2
{% from "kubernetes/map.jinja" import kubernetes with context %}

{% set component = 'kube-scheduler' %}
{% set component_bin_path = kubernetes.install_dir + '/scheduler' %}

{% from "kubernetes/vars.jinja" import
    kubernetes_apiserver_url,
    kubernetes_etc_dir,
    component_kubeconfig
with context %}

[Unit]
Description=Kubernetes Scheduler
Documentation=https://github.com/kubernetes/kubernetes
After=kube-apiserver.service

[Service]
User=root
ExecStart={{ component_bin_path }} \
    --kubeconfig={{ component_kubeconfig }} \
    --master={{ kubernetes_apiserver_url }} \
    --logtostderr=true
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
